# scripts/com.localnetinspector.agent.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.localnetinspector.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/yostrovs/Documents/New\ project/localnet-inspector && source .venv/bin/activate && localnet-inspector</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/localnetinspector.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/localnetinspector.err.log</string>
</dict>
</plist>
